# assets/supervisor/apache.conf
[program:apache]
command=/usr/sbin/apache2ctl -D FOREGROUND
environment=APACHE_RUN_USER="www-data",APACHE_RUN_GROUP="www-data",APACHE_LOG_DIR="/var/log/apache2"
autostart=true
autorestart=true
redirect_stderr=true
